# files/tuya_esp.init
#!/bin/sh /etc/rc.common
USE_PROCD=1
START=80
STOP=90
PROGRAM="/usr/bin/tuya_esp"
start_service() { 
    procd_open_instance tuya_esp
    procd_set_param command $PROGRAM
    procd_append_param command &
    procd_close_instance
}

service_triggers() {
    procd_add_reload_trigger "tuya_esp"
}

reload_service() {  
    stop
    start
}
